# systemd/pantry-scanner.service
[Unit]
Description=Pantry Barcode Scanner Service
After=network.target

[Service]
Type=simple
User=pantry
Group=pantry
ExecStart=/usr/bin/python3 /opt/shelf_manager_backend/scanner_service.py
Restart=on-failure
Environment="PYTHONUNBUFFERED=1"
EnvironmentFile=/opt/shelf_manager_backend/.env

[Install]
WantedBy=multi-user.target