# systemd/pantry-scanner.service
[Unit]
Description=Pantry Barcode Scanner Service
After=network.target

[Service]
Type=simple
User=pantry
Group=pantry
ExecStart=/usr/bin/python3 /home/pantry/shelf_manager_backend/pantry/scanner_service.py
Restart=on-failure
Environment="PYTHONUNBUFFERED=1"
EnvironmentFile=/home/pantry/shelf_manager_backend/pantry/.env

[Install]
WantedBy=multi-user.target